# root/etc/s6-overlay/s6-rc.d/service-autoscan/run
#!/command/with-contenv bash

umask "${UMASK}"

exec s6-setuidgid autoscan python3 /usr/src/app/src/scan.py ${AUTOSCAN_COMMAND}
